# crontab/test_cronjob.sh
#!/bin/bash

echo "$(date) : Test cron" >> cron_test.log
